# config/bf3-security.cron
# BlueField-3 security posture update for osquery
# Updates the bf3_security table every 5 minutes
*/5 * * * * root /usr/local/bin/bf3_security_update >/dev/null 2>&1
